# warehouse/sql-tunnels/driver/ssh/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D