# scripts/rudder.conf
# upstart configuration

description "Rudder server"
author "rudder"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /home/ubuntu/rudder-server/rudder-server
